# scripts/telegram-calendar-bot.service
[Unit]
Description=Telegram Google Calendar Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/telegram-calendar-bot
EnvironmentFile=/opt/telegram-calendar-bot/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
